# target/lemote/files/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/etc/init.d/rcK
tty::respawn:/sbin/getty -i -L tty 115200 vt100
tty1::respawn:/sbin/getty -i -L tty1 115200 vt100
