# package/base-files/extra/etc/network/if-up.d/01-bridge
#!/bin/sh

[ "${IFACE%%[0-9]*}" = "br" ] && ip link set up dev $IFACE
