# package/fwupdate/files/fwupdate.service
[Unit]
Description=Firmware Update Validation
After=dropbear.service

[Service]
ExecStart=/usr/sbin/fwvalidate

[Install]
WantedBy=multi-user.target
