# package/base-files/src/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/etc/init.d/rcK
